# linkbot-hub-announcer.service
# This file was autogenerated by node-deb {{ node_deb_version }}

[Unit]
Description={{ node_deb_package_description }}
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/usr/share/{{ node_deb_package_name }}
ExecStart=/usr/share/{{ node_deb_package_name }}/bin/{{ node_deb_executable_name }}
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10s
User=nobody

[Install]
WantedBy=multi-user.target
